# packaging/systemd/cns-matrix-web.service.in
# Serves the CNS matrix dashboard and control API on port 8888.
[Unit]
Description=CNS matrix dashboard web server
After=network.target

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_GROUP__
WorkingDirectory=/var/lib/cns-matrix/www
ExecStart=/usr/bin/python3 __REPO_ROOT__/tools/cns_matrix_web.py --repo-root __REPO_ROOT__ --base-dir /var/lib/cns-matrix --bind 0.0.0.0 --port 8888
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
